# monitored_environment/conf/firewall/rc.local
#!/bin/sh -e

/sbin/iptables -t nat -A POSTROUTING -o eth0 -j MASQUERADE

exit 0
